# quadlet-resource/tasks/container/template.container
[Container]
ContainerName=%name%
Descritpion=%description%
Image=%image%
PublishPort=%port%

[Install]
WantedBy=%targets%

